# linstor-gateway.service
[Unit]
Description=LINSTOR Gateway
After=network.target

[Service]
ExecStart=/usr/sbin/linstor-gateway server

[Install]
WantedBy=multi-user.target
